# cookbooks/cfrdevops/files/default/root.crontab
PATH=/sbin:/bin
*/2 * * * * service memcached restart && service memcached2 restart && service memcached3 restart >/dev/null 2&>1 

